# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:migrations]
command=python -m reflexio_ext.scripts.run_all_migrations --continue-on-error
directory=/app
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fastapi]
command=python -m uvicorn reflexio_ext.server.api:app --host 0.0.0.0 --port 8081
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nextjs]
command=npm run start --prefix /app/website
environment=NODE_OPTIONS="--max-old-space-size=512"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:docs]
command=npm run start --prefix /app/public_docs -- -p 8082
environment=NODE_OPTIONS="--max-old-space-size=256"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=25
